# doorlock-web.service
[Unit]
Description=Door Lock Web Interface
After=network.target

[Service]
Type=simple
User=guillaume
Group=guillaume
WorkingDirectory=/home/guillaume/serrure_automatique
ExecStart=/usr/bin/python3 /home/guillaume/serrure_automatique/web_app.py --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ensure camera and GPIO access for web interface
SupplementaryGroups=gpio video

[Install]
WantedBy=multi-user.target
